# v1/opt/auto-update/fleet/update-check.timer
[Unit]
Description=Runs CoreOS update check every 30m

[Timer]
OnUnitActiveSec=30m

[X-Fleet]
Global=true
