# v1/fleet/update-os.timer
[Unit]
Description=Runs CoreOS update check every 10m

[Timer]
OnBootSec=10min

OnUnitActiveSec=10min

[X-Fleet]
Global=true
